# services/backup.service
[Unit]
Description=Runs backup everytime the laptop boots

[Service]
ExecStart=/usr/bin/backup

[Install]
WantedBy=multi-user.target 
